# services/bme280/bme280.service
[Unit]
Description=Arkadia BME280 Climate Sensor Service
Documentation=https://github.com/malonge/Arkadia
After=mosquitto.service
Requires=mosquitto.service
# Rate-limit restarts: no more than 5 restarts in 60 s.
# Must be in [Unit] (not [Service]) on systemd 229+; in [Service] the keys
# are silently ignored as "Unknown key name".
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# Change 'pi' to your actual username.
User=pi
Group=i2c

# /etc/home-monitor.env must exist and define ARKADIA_ROOT — the absolute
# path to this repository on the target machine.  Example:
#   ARKADIA_ROOT=/home/pi/Projects/Arkadia
# systemd expands ${ARKADIA_ROOT} in ExecStart from this file.
EnvironmentFile=/etc/home-monitor.env

# Full paths are used here because WorkingDirectory= does not expand
# environment variables.  sys.path in main.py ensures sensor.py is found
# regardless of the working directory.
#
# Note: 'systemd-analyze verify' will report a false positive here because
# that tool does not load EnvironmentFile during static analysis.  At
# runtime the systemd daemon DOES expand ${ARKADIA_ROOT} from the env file
# before executing the command.
ExecStart=${ARKADIA_ROOT}/services/bme280/.venv/bin/python \
          ${ARKADIA_ROOT}/services/bme280/main.py

Restart=on-failure
RestartSec=5
# Allow up to 15 s for a clean shutdown (covers the worst-case poll cycle).
TimeoutStopSec=15

# Log via journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bme280

# Hardening
NoNewPrivileges=true
PrivateTmp=true
# 'full' makes /usr, /boot, /etc read-only; /dev (I2C) remains accessible.
ProtectSystem=full

[Install]
WantedBy=multi-user.target
